# config/autostart/unclutter.desktop
[Desktop Entry]
Type=Application
Name=unclutter
NoDisplay=true
Comment=Hides the cursor on inactivity
Exec=unclutter --timeout 1
